# etc/services.d/generate_thumbs/run
#!/usr/bin/execlineb -P
with-contenv
fdmove -c 2 1
cd /app
s6-envuidgid django
s6-applyuidgid -U
bash -ex -c "while nice -n 10 ./manage.py generate_thumbs; do sleep 7200; done"