# etc/services.d/gunicorn/run
#!/usr/bin/execlineb -P
with-contenv
fdmove -c 2 1
cd /app
/usr/local/bin/gunicorn -c ./gunicorn.conf uqam.wsgi